# terraform/crontab
0 9 * * 1-5 terraform -chdir=terraform -auto-approve apply -
